# crontab
* * * * * scripts/path_of_the_script_file
